# run.crontab
00 06,18 * * * /usr/bin/python /home/pi/projects/sprinkler/run_sprinkler.py