# run.crontab.sample
 00 06,18 * * * /usr/bin/python2.7 /home/pi/rpi_sprinkler/run_sprinkler.py
