# run.crontab.sample
 00 06,18 * * * /usr/bin/python /home/pi/rpi_sprinkler/run_sprinkler.py
